# tdt/cvs/cdk/static/common/etc/inittab

id:3:initdefault:

si::sysinit:/etc/init.d/rcS

te:2:respawn:/sbin/getty -L ttyAS0 115200 vt102
st:3:respawn:/etc/init.d/start

ha:0:wait:/etc/init.d/halt
re:6:wait:/etc/init.d/reboot
